# systemd/lumen-occupancy-daily.timer
[Unit]
Description=Run Lumen occupancy report at 09:00

[Timer]
OnCalendar=*-*-* 09:00:00
Unit=lumen-occupancy-daily.service

[Install]
WantedBy=timers.target
